# assets/contrib/deb/prerm
systemctl stop medveddb
systemctl disable medveddb
